# configs/etc/systemd/system/php8.2-fpm.service.d/override.conf
[Service]
Environment="GITHUB_WEBHOOK_SECRET=dcMbplDLGNg0sSv43uJZ5Ws5nW5KFxLZReYCpC2oCGw="